# 备份/supervisor_templates/member_bot3.conf
# Supervisor配置 - 会员机器人3

[program:member_bot3]
command=/home/ubuntu/telegram_bots/bybot3.0/venv/bin/python /home/ubuntu/telegram_bots/bybot3.0/lsjmain.py
directory=/home/ubuntu/telegram_bots/bybot3.0
environment=ENV_FILE=".env.member_bot3"
user=ubuntu
autostart=true
autorestart=true
startretries=10
stderr_logfile=/var/log/member_bot3.err.log
stdout_logfile=/var/log/member_bot3.out.log
redirect_stderr=true
stopwaitsecs=30
killasgroup=true
stopasgroup=true








