# supervisor_templates/member_bot2.conf
# Supervisor配置 - 会员机器人2

[program:member_bot2]
command=/home/ubuntu/telegram_bots/bybot3.0/venv/bin/python /home/ubuntu/telegram_bots/bybot3.0/lsjmain.py
directory=/home/ubuntu/telegram_bots/bybot3.0
environment=ENV_FILE=".env.member_bot2"
user=ubuntu
autostart=true
autorestart=true
startretries=10
stderr_logfile=/var/log/member_bot2.err.log
stdout_logfile=/var/log/member_bot2.out.log
redirect_stderr=true
stopwaitsecs=30
killasgroup=true
stopasgroup=true








